# install/gun-configs/gunicorn.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=127.0.0.1:5003

[Install]
WantedBy=sockets.target
